# states/apps/lavalink-server/files/lavalink.service.jinja
[Unit]
Description=Lavalink Server
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/lavalink
# Ensure the log directory exists and has proper ownership.
ExecStartPre=/bin/bash -c 'sudo mkdir -p /var/log/lavalink && sudo chown ubuntu:ubuntu /var/log/lavalink'
# Start Lavalink using Java. Adjust the java executable path if needed.
ExecStart=/bin/bash -c 'exec /usr/bin/java -jar /opt/lavalink/lavalink.jar >> /var/log/lavalink/lavalink_$(date +%%Y%%m%%d).log 2>&1'
# Only set the essential environment variables using pillar data.
Environment="SERVER_HOST={{ salt['pillar.get']('config', {}).get('LAVALINK', {}).get('LAVALINK_HOST', 'http://127.0.0.1') }}"
Environment="SERVER_PORT={{ salt['pillar.get']('config', {}).get('LAVALINK', {}).get('LAVALINK_PORT', 2333) }}"
Environment="LAVALINK_SERVER_PASSWORD={{ salt['pillar.get']('config', {}).get('LAVALINK', {}).get('LAVALINK_PASSWORD', 'youshallnotpass') }}"
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
